# ~/Library/LaunchAgents/com.revyl.cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.your_company.cloudflared</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>your_company-bot</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string><HOME>/Library/Logs/cloudflared.err.log</string>
    <key>StandardOutPath</key>
    <string><HOME>/Library/Logs/cloudflared.out.log</string>
    <key>WorkingDirectory</key>
    <string><HOME></string>
</dict>
</plist> 